# rootfs/etc/services.d/apache/run
#!/usr/bin/with-contenv bash

/usr/local/bin/apache2-foreground
